# com.johnserious.mbo-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.johnserious.mbo-mcp</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/johnserious/MBO/scripts/mbo-start.sh</string>
        <string>--mode=dev</string>
        <string>--agent=operator</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/johnserious/MBO</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/johnserious</string>
        <key>NVM_DIR</key>
        <string>/Users/johnserious/.nvm</string>
        <key>MBO_PORT</key>
        <string>3737</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>3</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>ExitTimeout</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/Users/johnserious/MBO/.dev/logs/mbo-mcp.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/johnserious/MBO/.dev/logs/mbo-mcp-error.log</string>
</dict>
</plist>
