# scripts/com.mbo.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mbo.mcp</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/johnserious/MBO/scripts/mbo-start.sh</string>
        <string>--mode=dev</string>
        <string>--agent=launchd</string>
    </array>

    <!-- Working directory must be repo root for relative path resolution -->
    <key>WorkingDirectory</key>
    <string>/Users/johnserious/MBO</string>

    <!-- nvm node must be on PATH — launchd does not source shell profiles -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/johnserious/.nvm/versions/node/v24.11.1/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NVM_DIR</key>
        <string>/Users/johnserious/.nvm</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/johnserious/MBO/.dev/logs/mcp.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/johnserious/MBO/.dev/logs/mcp.err</string>

    <!-- Throttle respawn — prevents tight crash loop -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
